# systemd/librespot-check.service
[Unit]
Description=Check if Spotify Connect device is discoverable

[Service]
Type=oneshot
ExecStart=/home/pi/librespot-check.sh
